# .ebextensions/cronjob.cron
45 04 28 */1 * source /opt/python/current/env && cd /opt/python/current/app/edgar && /opt/python/run/venv/bin/python3 /opt/python/current/app/manage.py runcrons >> /var/log/cronjobs.log