# .ebextensions/cronjob.cron
* * * * * root /opt/python/run/venv/bin/python3 /opt/python/fedgehundapi/manage.py runcrons >> /var/log/cronjobs.log