# .ebextensions/cronjob.cron
* * * * * /bin/date >> /tmp/cron_output.log 2>&1